# .circle/supervisor_ph.conf
[program:PH_NAME]
stopasgroup=true
user=deploy-user
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/PH_NAME.log
stderr_logfile=/var/log/supervisor/PH_NAME.log
redirect_stderr=true